# user/default/files/etc/hotplug.d/iface/21-domain
#!/bin/bash
echo "address=/o.to/`uci show network.lan.ipaddr | awk -F"'" '{print $2}'`" > /tmp/dnsmasq.d/21-domain.conf
/etc/init.d/dnsmasq restart
